# php-cli/7.2-supv-cron/start/conf/supervisor.d/crontab.ini
[program:crontab]
process_name=%(program_name)s
command=/usr/sbin/crond -f
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true